# example/rootfs/etc/s6-overlay/s6-rc.d/example1/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# Minimal hello service
bashio::log.info "Starting Hello service..."
exec /usr/bin/hello.sh
